# crontab.cron
0 8 * * 1 echo Happy Monday Morning
* * * * * echo One Minute Passed > /dev/ttyl
30 14 1 * 1 mail users % Jan Meeting At 3pm
